# crontab
* * * * * /home/parallels/joker.sh
* * * * * /home/parallels/joker.sh
